# templates/init/carbon-relay.conf
description "Graphite carbon-relay server"
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/graphite/bin/carbon-relay.py start
